# ssserver.conf
description "ShadowSocks Server"

start on net-device-up
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/local/bin/ssserver -c /etc/shadowsocks.json

